# crontab
* * * * * ~/bin/ping_gw.sh
* * * * * ~/bin/get_default_gw.sh
